# util/70-snap.barrier.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_barrier_barrier"
# opengl
KERNEL=="nvmap", TAG+="snap_barrier_barrier"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_barrier_barrier"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_barrier_barrier"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_barrier_barrier"
# opengl
KERNEL=="vchiq", TAG+="snap_barrier_barrier"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_barrier_barrier"
TAG=="snap_barrier_barrier", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_barrier_barrier $devpath $major:$minor"
# opengl
KERNEL=="nvhost-*", TAG+="snap_barrier_barrierc"
# opengl
KERNEL=="nvmap", TAG+="snap_barrier_barrierc"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_barrier_barrierc"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_barrier_barrierc"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_barrier_barrierc"
# opengl
KERNEL=="vchiq", TAG+="snap_barrier_barrierc"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_barrier_barrierc"
TAG=="snap_barrier_barrierc", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_barrier_barrierc $devpath $major:$minor"
# opengl
KERNEL=="nvhost-*", TAG+="snap_barrier_barriers"
# opengl
KERNEL=="nvmap", TAG+="snap_barrier_barriers"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_barrier_barriers"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_barrier_barriers"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_barrier_barriers"
# opengl
KERNEL=="vchiq", TAG+="snap_barrier_barriers"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_barrier_barriers"
TAG=="snap_barrier_barriers", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_barrier_barriers $devpath $major:$minor"
